# gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/home/ubuntu/books_reviews/jktech_app/jktech_app.sock

[Install]
WantedBy=sockets.target
